# contrib/init/GuldenD.conf
description "Gulden Core Daemon"

start on runlevel [2345]
stop on starting rc RUNLEVEL=[016]

env GULDEND_BIN="/usr/bin/GuldenD"
env GULDEND_USER="gulden"
env GULDEND_GROUP="gulden"
env GULDEND_PIDDIR="/var/run/GuldenD"
# upstart can't handle variables constructed with other variables
env GULDEND_PIDFILE="/var/run/GuldenD/GuldenD.pid"
env GULDEND_CONFIGFILE="/etc/gulden/gulden.conf"
env GULDEND_DATADIR="/var/lib/GuldenD"

expect fork

respawn
respawn limit 5 120
kill timeout 60

pre-start script
    # this will catch non-existent config files
    # GuldenD will check and exit with this very warning, but it can do so
    # long after forking, leaving upstart to think everything started fine.
    # since this is a commonly encountered case on install, just check and
    # warn here.
    if ! grep -qs '^rpcpassword=' "$GULDEND_CONFIGFILE" ; then
        echo "ERROR: You must set a secure rpcpassword to run GuldenD."
        echo "The setting must appear in $GULDEND_CONFIGFILE"
        echo
        echo "This password is security critical to securing wallets "
        echo "and must not be the same as the rpcuser setting."
        echo "You can generate a suitable random password using the following"
        echo "command from the shell:"
        echo
        echo "bash -c 'tr -dc a-zA-Z0-9 < /dev/urandom | head -c32 && echo'"
        echo
        echo "It is also recommended that you also set alertnotify so you are "
        echo "notified of problems:"
        echo
        echo "ie: alertnotify=echo %%s | mail -s \"Gulden Alert\"" \
            "admin@foo.com"
        echo
        exit 1
    fi

    mkdir -p "$GULDEND_PIDDIR"
    chmod 0755 "$GULDEND_PIDDIR"
    chown $GULDEND_USER:$GULDEND_GROUP "$GULDEND_PIDDIR"
    chown $GULDEND_USER:$GULDEND_GROUP "$GULDEND_CONFIGFILE"
    chmod 0660 "$GULDEND_CONFIGFILE"
end script

exec start-stop-daemon \
    --start \
    --pidfile "$GULDEND_PIDFILE" \
    --chuid $GULDEND_USER:$GULDEND_GROUP \
    --exec "$GULDEND_BIN" \
    -- \
    -pid="$GULDEND_PIDFILE" \
    -conf="$GULDEND_CONFIGFILE" \
    -datadir="$GULDEND_DATADIR" \
    -disablewallet \
    -daemon

